# files_rpi/etc/rc.local
#!/bin/sh -e

ln -sf /bin/bash /bin/sh
exit 0
