# files_rpi/tmp/home/tty2rpi/.config/systemd/user/tty2rpi-plocate.timer
[Unit]
Description=Update media library every minute

[Timer]
OnBootSec=10s
OnCalendar=*-*-* *:*:00
Persistent=true
AccuracySec=1ms

[Install]
WantedBy=timers.target
