# files_rpi/tmp/home/tty2rpi/.config/systemd/user/tty2rpi-screensaver.timer
[Unit]
Description=Run the tty2rpi screensaver

[Timer]
OnActiveSec=120
OnUnitActiveSec=10
AccuracySec=1ms

[Install]
WantedBy=timers.target

