# files_rpi/tmp/home/tty2rpi/.config/systemd/user/tty2rpi-update.timer
[Unit]
Description=Weekly update of tty2rpi

[Timer]
OnCalendar=Fri 14:00
RandomizedDelaySec=7200
Persistent=true
AccuracySec=1ms

[Install]
WantedBy=timers.target

