# files_rpi/tmp/home/tty2rpi/.config/systemd/user/tty2rpi-plocate.service
[Unit]
Description=Update media library

[Service]
ExecStart=%h/tty2rpi-scripts/tty2rpi-plocate.sh
