# files_rpi/tmp/home/tty2rpi/.config/systemd/user/tty2rpi-screensaver.service
[Unit]
Description=tty2rpi screensaver

[Service]
ExecStart=%h/tty2rpi-scripts/tty2rpi-clock.sh
